# root/etc/services.d/youtube-dl/run
#!/usr/bin/with-contenv bash

sleep 60s
exec s6-setuidgid abc /app/youtube-dl/youtube-dl.sh
