# root/etc/services.d/youtube-dl/run
#!/usr/bin/with-contenv bash

exec s6-setuidgid nobody /home/nobody/youtube-dl.sh
