# root/etc/services.d/youtube-dl/run
#!/usr/bin/with-contenv bash

sleep 1s
umask "$UMASK"
fdmove -c 2 1 exec s6-setuidgid abc /app/youtube-dl/youtube-dl.sh
