# root/etc/services.d/youtube-dl/run
#!/usr/bin/with-contenv bash

exec s6-setuidgid abc /app/youtube-dl/youtube-dl.bash
