# root/etc/services.d/youtube-dl/run
#!/usr/bin/with-contenv sh

exec s6-setuidgid abc:abc /app/youtube-dl/youtube-dl.sh
