# root/etc/supervisor/conf.d/terminate.conf
[program:terminate]
command=/etc/supervisor/terminate.sh
autorestart=false
autostart=false
