# etc/services.d/youtube-dl/run
#!/usr/bin/with-contenv sh

apk update > /dev/null
apk upgrade > /dev/null

if [ -f "/config/archive.txt" ]
then
  if [ ! -f "/config/dateafter.txt" ]
  then
    date "+%Y%m%d" > "/config/dateafter.txt"
  fi
else
  touch "/config/archive.txt"
  date "+%Y%m%d" > "/config/dateafter.txt"
fi

youtube-dl \
  --format "bestvideo[height<=$youtubedl_quality]+bestaudio[acodec!=opus]" \
  --config-location "/config/args.conf" \
  --download-archive "/config/archive.txt" \
  --dateafter "$(cat "/config/dateafter.txt")" \
  --batch-file "/config/channels.txt"

date
echo "Waiting $youtubedl_interval.."
sleep $youtubedl_interval
date
