# etc/services.d/youtube-dl/run
#!/usr/bin/with-contenv sh

apk update > /dev/null
apk upgrade > /dev/null

if [ ! -f "/config/archive.txt" ]
then
  if [ "$youtubedl_youtube_login" = "yes" ]
  then
    youtube-dl \
      --format "bestvideo[height<=$youtubedl_quality]+bestaudio[acodec!=opus]" \
      --config-location "/config/args.conf" \
      --username $youtubedl_youtube_username \
      --password $youtubedl_youtube_password \
      --playlist-start 2 \
      --get-id
  else
    youtube-dl \
      --format "bestvideo[height<=$youtubedl_quality]+bestaudio[acodec!=opus]" \
      --config-location "/config/args.conf" \
      --batch-file "/config/channels.txt" \
      --playlist-start 2 \
      --get-id
  fi
fi

if [ "$youtubedl_youtube_login" = "yes" ]
then
  youtube-dl \
    --format "bestvideo[height<=$youtubedl_quality]+bestaudio[acodec!=opus]" \
    --config-location "/config/args.conf" \
    --username $youtubedl_youtube_username \
    --password $youtubedl_youtube_password
else
  youtube-dl \
    --format "bestvideo[height<=$youtubedl_quality]+bestaudio[acodec!=opus]" \
    --config-location "/config/args.conf" \
    --batch-file "/config/channels.txt"
fi
date
echo "sleeping $youtubedl_interval.."
sleep $youtubedl_interval
