# etc/services.d/youtube-dl/run
#!/usr/bin/with-contenv sh

echo "checking for updates.."
apk update
apk upgrade

if [ -f "/config/archive.txt" ]
then
  if [ ! -f "/config/dateafter.txt" ]
  then
    echo "'dateafter.txt' not found."
    echo "creating 'dateafter.txt'."
    date --date "-1 days" "+%Y%m%d" > "/config/dateafter.txt"
  fi
else
  echo "'archive.txt' not found."
  echo "creating 'archive.txt'."
  touch "/config/archive.txt"
  echo "creating 'dateafter.txt'."
  date --date "-1 days" "+%Y%m%d" > "/config/dateafter.txt"
fi

youtube-dl \
  --format "bestvideo[height<=$youtubedl_quality]+bestaudio[acodec!=opus]" \
  --config-location "/config/args.conf" \
  --download-archive "/config/archive.txt" \
  --dateafter "$(cat "/config/dateafter.txt")" \
  --batch-file "/config/channels.txt"

date
echo "waiting $youtubedl_interval.."
sleep $youtubedl_interval
date
