# root/etc/services.d/updater/run
#!/usr/bin/with-contenv bash

exec /app/youtube-dl/updater.sh
