# root/etc/services.d/updater/run
#!/usr/bin/with-contenv bash

exec s6-envuidgid root /app/youtube-dl/updater.sh
